# ansible/files/pi-disable-hdmi.service
[Unit]
Description=PI disable HDMI (save power)

[Service]
Type=oneshot
ExecStart=/usr/bin/tvservice -o
